# supervisord.conf
[program:celery]
command=celery -A app.celery worker
[program:uwsgi]
command=uwsgi --ini ./app.ini